# deploy/gunicorn.socket
[Unit]
Description=gunicorn socket for Sana Academy

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
