# deploy/gunicorn.socket
[Unit]
Description=Gunicorn socket for SNA Academy

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
